# pkg-debian/usr/bin/alarm/alarm.service
[Unit]
Description=Launching Alarm Server with systemd, respawn
After=network.target

[Service]
ExecStart=/usr/bin/alarm/alarm-server
Restart=always

[Install]
WantedBy=multi-user.target